# scripts/inspiron/restic-vps-backup.service
[Unit]
Description=Trigger VPS Restic backup when external HDD is connected
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-trigger-vps-backup.sh
